# deploy/ai-stock-app.service
[Unit]
Description=AI Stock App Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ai-stock-app/backend
Environment="PATH=/opt/ai-stock-app/venv/bin"
ExecStart=/opt/ai-stock-app/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
